# root/etc/s6-overlay/s6-rc.d/svc-ckpool/run
#!/usr/bin/with-contenv bash

exec /app/bin/ckpool -c /config/ckpool.conf
